# Downloads/pharmacy-main(1)/pharmacy-main/pharmacy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:php-fpm]
command=php-fpm --nodaemonize
autorestart=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=20
